# roles/counting-panda/files/counting-panda.conf
[program:counting-panda]
directory=/tmp/counting-panda
command=nodejs counting-panda.js
stdout_logfile=/var/log/counting-panda.log
stderr_logfile=/var/log/error_counting-panda.log

