# roles/static-panda/files/static-panda.conf
[program:static-panda]
directory=/tmp/static-panda
command=nodejs static-panda.js
stdout_logfile=/var/log/static-panda.log
stderr_logfile=/var/log/error_static-panda.log

